# lxde/lxde-dev/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/usr/bin/with-contenv bash
exec /var/www/html/utils/novnc_proxy --vnc localhost:5900 --listen 6080
