# lxde/lxde-11/etc/s6-overlay/s6-rc.d/svc-xvfb/run
#!/bin/sh
export DISPLAY=":1" HOME="/root" USER="root"
exec /usr/local/bin/xvfb.sh
