# src/alpine/install/thunderbird/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0

[program:x11]
priority=0
user=thunderbird
command=/usr/bin/Xvnc -desktop "Thunderbird" -localhost -rfbport 5900 \
-SecurityTypes VncAuth -PasswordFile /thunderbird/.vnc/passwd -AlwaysShared \
-AcceptKeyEvents -AcceptPointerEvents -AcceptSetDesktopSize -SendCutText \
-AcceptCutText :0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true true

[program:novnc]
priority=0
command=/usr/local/novnc/utils/novnc_proxy --vnc localhost:5900
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:openbox]
priority=1
command=/usr/bin/openbox
environment=DISPLAY=:0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:thunderbird]
priority=1
environment=DISPLAY=:0
command=/usr/bin/thunderbird
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
