# lxde/lxde-11/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/bin/sh
exec /home/ubuntu/utils/novnc_proxy --vnc localhost:5900 --listen 6080
